# modules/slave/config/files/token-refresher.timer
[Unit]
Description=Refresh the Token ever 8 hours

[Timer]
OnBootSec=0
OnUnitActiveSec=8h

[Install]
WantedBy=timers.target
